# deploy/systemd/auto_review.service
[Unit]
Description=auto_review (Forgejo PR-review bot)
Documentation=https://github.com/jwilger/auto_review
# Forgejo needs to be reachable; in particular /readyz probes it.
# We don't gate on a specific Forgejo unit name — operators using
# Forgejo on the same host can add `Wants=forgejo.service` and
# `After=forgejo.service` in a drop-in.
After=network-online.target
Wants=network-online.target
# Burst guard: > 5 restarts in 60s = give up so we don't mask a
# permanent misconfiguration. systemd parses these in [Unit]; the
# matching `Restart=` directive lives in [Service] below.
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=exec
# Operator-managed dedicated user. Create with:
#   useradd --system --no-create-home --shell /usr/sbin/nologin auto_review
User=auto_review
Group=auto_review

# All credentials and tunables live in this file. See
# `deploy/systemd/auto_review.env.example` for the full list.
EnvironmentFile=/etc/auto_review/auto_review.env

# The binary is built with `cargo build --release -p ar-cli`;
# operators install it under /usr/local/bin or symlink from a pinned
# location. Adjust for your layout.
ExecStart=/usr/local/bin/auto-review gateway

# Workspace tmpfs for the orchestrator's clones. Smaller systems
# can drop this and let clones go to /tmp; bigger ones with
# heavy traffic should size this for their largest expected PR.
RuntimeDirectory=auto_review
RuntimeDirectoryMode=0700

# Persistent state for SQLite-backed learnings, review history, vector context,
# and webhook dedup when the env file uses /var/lib/auto_review paths.
# RuntimeDirectory is wiped on restart; StateDirectory survives.
StateDirectory=auto_review
StateDirectoryMode=0700
ReadWritePaths=/var/lib/auto_review

# Restart on transient failures (LLM hiccup, Forgejo blip).
# The matching StartLimit* directives are in [Unit] above
# (systemd insists they live there).
Restart=on-failure
RestartSec=5s

# Graceful-shutdown semantics. The gateway listens for SIGTERM
# and drains its HTTP listener before exiting; give it room to
# finish in-flight requests rather than hard-killing at the
# default 90s. Reviews already spawned via tokio::spawn are
# best-effort: if they're still running at the timeout, systemd
# escalates to SIGKILL.
KillSignal=SIGTERM
TimeoutStopSec=30s

# Hardening for the bot process itself. Runtime linter execution was retired;
# these directives are now the service-level guardrail if a future bug lets the
# gateway process misbehave.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
PrivateTmp=true
PrivateDevices=true
PrivateUsers=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallFilter=@system-service
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits tuned for a single-tenant self-host on a
# 4-core / 8 GiB box reviewing ~hundreds of PRs/day. Bump for
# a busier instance.
LimitNOFILE=4096
TasksMax=512

# Standard logging targets so `journalctl -u auto_review` works
# out of the box.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auto_review

[Install]
WantedBy=multi-user.target
